# anacrontab
MAILTO=felixc@felixcrux.com
SHELL=/bin/sh
PATH=$HOME/bin:/usr/local/bin:/usr/local/sbin:/usr/bin:/usr/sbin:/bin:/sbin

1         10   backup-daily    nice $HOME/bin/backup.sh daily
7         20   backup-weekly   nice $HOME/bin/backup.sh weekly
@monthly  30   backup-monthly  nice $HOME/bin/backup.sh monthly
