# crontab
@reboot /usr/sbin/anacron -s -S $HOME/.local/state/anacronspool -t $HOME/cfg/anacrontab
@hourly /usr/sbin/anacron -s -S $HOME/.local/state/anacronspool -t $HOME/cfg/anacrontab
