# crontab
MAILTO=felixc@felixcrux.com
RANDOM_DELAY=5

@reboot /usr/sbin/anacron -s -S $HOME/.anacronspool -t $HOME/.anacrontab
@hourly /usr/sbin/anacron -s -S $HOME/.anacronspool -t $HOME/.anacrontab
