# crontab
@reboot /usr/sbin/anacron -s -S $HOME/.anacronspool -t $HOME/.anacrontab
@hourly /usr/sbin/anacron -s -S $HOME/.anacronspool -t $HOME/.anacrontab
